# othersrc/lib-systemd-system/wslg-xwayland.service
[Unit]
Requires=wslg-xwayland.socket
After=wslg-xwayland.socket

[Service]
ExecStart=/lib/systemd/systemd-socket-proxyd /mnt/wslg/.X11-unix/X0
